# sys-files/nintel05/sys-files.list
# one file/dir name on one line with absolute path

/etc/environment

/home/codeman/.zshrc
/home/codeman/.zsh_history

/home/codeman/.bashrc

/home/codeman/.config/i3/config

/etc/systemd/system/ad-i3lock.service
